# data/user/steamos-manager.service
[Unit]
Description=SteamOS Manager Daemon
After=gamescope.service

[Service]
Type=dbus
BusName=com.steampowered.SteamOSManager1
Environment=RUST_LOG='INFO'
ExecStart=/usr/lib/steamos-manager
Restart=on-failure
RestartSec=1
