# data/steamos-manager.service
[Unit]
Description=SteamOS Manager Daemon

[Service]
Type=dbus
BusName=com.steampowered.SteamOSManager1
ExecStart=/usr/lib/steamos-manager
